# main.cron
27 3 * * * /home/de10/fzlv/project/main.py
